# overlay/etc/uci-defaults/91-hx-network-init
#!/bin/sh
set -e

lan_ip_address="192.168.18.1"

exec >/tmp/hx-uci-20-lan.log 2>&1

# Configure LAN (IP + netmask 必须设置！)
if [ -n "$lan_ip_address" ]; then
  uci set network.lan.ipaddr="$lan_ip_address"
  uci set network.lan.netmask='255.255.255.0'
  uci commit network
fi

echo "[OK] lan ip done"
exit 0
