# overlay/etc/uci-defaults/97-hx-openclash-core-init
#!/bin/sh
set -eu

LOG="/tmp/hx-uci-97-openclash.log"
TAG="hx-openclash-init"

echo "[$(date -Is)] uci-defaults start" >> "$LOG"

has_openclash() {
  [ -f /usr/lib/lua/luci/controller/openclash.lua ] && return 0
  [ -f /etc/config/openclash ] && return 0
  [ -d /usr/share/openclash ] && return 0
  return 1
}

if ! has_openclash; then
  echo "[$(date -Is)] openclash not installed, skip" >> "$LOG"
  exit 0
fi

# run 脚本必须存在
if [ ! -x /usr/bin/hx-openclash-init-run ]; then
  echo "[$(date -Is)] /usr/bin/hx-openclash-init-run missing, skip" >> "$LOG"
  exit 0
fi

# 装 init.d（如果你已经打进固件可省略）
if [ ! -x /etc/init.d/hx-openclash-init ]; then
  cat > /etc/init.d/hx-openclash-init <<'EOF'
#!/bin/sh /etc/rc.common
START=99
USE_PROCD=1
start_service() {
  procd_open_instance
  # 延迟，等待 hx-storage 先把 /mnt/custom 挂好
  procd_set_param command /bin/sh -c 'sleep 30; /usr/bin/hx-openclash-init-run'
  procd_set_param respawn 0 0 0
  procd_close_instance
}
EOF
  chmod +x /etc/init.d/hx-openclash-init
fi

/etc/init.d/hx-openclash-init enable || true

echo "[$(date -Is)] enabled hx-openclash-init" >> "$LOG"
logger -t "$TAG" "uci-defaults enabled init.d (see $LOG)"
exit 0
