# overlay/etc/uci-defaults/99-hx-finish
#!/bin/sh
set -e

exec >/tmp/hx-uci-90-finish.log 2>&1

# Apply network/dhcp changes
/etc/init.d/network restart || true
/etc/init.d/dnsmasq restart || true

# # Apply wifi (if wifi exists)
wifi reload || wifi up || true

# # ✅ Ensure firewall is enabled and running
/etc/init.d/firewall enable || true
/etc/init.d/firewall restart || /etc/init.d/firewall start || true
# /etc/init.d/rpcd restart || true

# ---- Fix uhttpd slow access (remove invalid lua_prefix) ----
uci -q delete uhttpd.main.lua_prefix

# 放宽一点限制，避免首次访问慢
uci -q set uhttpd.main.max_requests='0'
uci -q set uhttpd.main.max_connections='100'

uci -q commit uhttpd
/etc/init.d/uhttpd restart || true


echo "[OK] finish done"
exit 0
  
   
 
